# docker/crontab
*/1 * * * * /usr/bin/node /var/www/api/cli.js push-task services.checkAll >/dev/null 2>&1

* * * * * root echo "Hello world"
