# docker/crontab
*/1 * * * * root /usr/local/bin/node /src/api/cli.js push-task services.checkAll >/dev/null 2>&1
